# systemd/mdcheck-ng.timer
[Unit]
Description=Periodically runs mdcheck-ng

[Timer]
OnCalendar=daily


[Install]
WantedBy=timers.target
